# launchd/com.maverick.tiingoloader.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maverick.tiingoloader.weekly</string>

    <key>WorkingDirectory</key>
    <string>/Users/overman/projects/maverick-mcp</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/overman/.local/bin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--env-file</string>
        <string>/Users/overman/projects/maverick-mcp/.env</string>
        <string>scripts/load_tiingo_data.py</string>
        <string>--sp500</string>
        <string>--years</string>
        <string>2</string>
        <string>--calculate-indicators</string>
        <string>--run-screening</string>
    </array>

    <!-- 每周六 北京时间 14:00 深度更新 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>6</integer>
        <key>Hour</key><integer>14</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardErrorPath</key>
    <string>/Users/overman/projects/maverick-mcp/logs/tiingo.weekly.log</string>
</dict>
</plist>
