# launchd/com.maverick.tiingoloader.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maverick.tiingoloader.daily</string>

    <key>WorkingDirectory</key>
    <string>/Users/overman/projects/maverick-mcp</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/overman/.local/bin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/overman/.local/bin/uv</string>
        <string>run</string>
        <string>--env-file</string>
        <string>/Users/overman/projects/maverick-mcp/.env</string>
        <string>scripts/load_tiingo_data.py</string>
        <string>--sp500</string>
        <string>--years</string>
        <!-- 约1个多月数据，适合每日增量更新 -->
        <string>0.1</string>
        <string>--run-screening</string>
    </array>

    <!-- 每天周一到周五 北京时间 07:00 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/overman/projects/maverick-mcp/logs/tiingo.daily.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/overman/projects/maverick-mcp/logs/tiingo.daily.err.log</string>
</dict>
</plist>
