# launchd/com.maverick.tiingoloader.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maverick.tiingoloader.weekly</string>

    <key>WorkingDirectory</key>
    <string>/Users/overman/Workspace/maverick-mcp</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
        <key>DATABASE_URL</key>
        <string>postgresql://postgres:ZCBsBZhHAMKPj9RTwHHn@localhost:5432/maverick_mcp</string>
        <key>TIINGO_API_TOKEN</key>
        <string>32f50ec7392ce4c52c664c1a61854388ff3e6e20</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>scripts/load_tiingo_data.py</string>
        <string>--sp500</string>
        <!-- 深度刷新，拉取过去2年的完整数据并强制重新计算技术指标 -->
        <string>--years</string>
        <string>2</string>
        <string>--calculate-indicators</string>
        <string>--run-screening</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <!-- 每周六 (Weekday=6) 的下午 14:00 (北京时间，对应美国深夜 2:00) -->
        <key>Weekday</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>14</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/overman/Workspace/maverick-mcp/logs/tiingo.weekly.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/overman/Workspace/maverick-mcp/logs/tiingo.weekly.err.log</string>
</dict>
</plist>
